# tado_optimizer/rootfs/etc/s6-overlay/s6-rc.d/init-tado_optimizer/run
#!/command/with-contenv bashio
# shellcheck shell=bash
# ==============================================================================
# Home Assistant Community Add-on: Example
# Example init script, runs before any other service
# ==============================================================================

bashio::log.info "Tado Optimizer init script running"
exec /run.sh
