# tado_heating_control/rootfs/etc/s6-overlay/s6-rc.d/init-tado_heating_control/up
/etc/s6-overlay/s6-rc.d/init-tado_heating_control/run
